# node-exporter/dockerinfo/dockerinfo.cron
*/5 * * * * root /usr/local/bin/dockerinfo.sh >>/var/log/dockerinfo.log
