# node-exporter/dockerinfo/dockerinfo.cron
*/1 * * * * root /usr/local/bin/dockerinfo.sh
